# project/config/sphinx/crontab
*/1 * * * * /Users/ProtasovSE/PycharmProjects/mlo_rest/mlo_rest/project/config/sphinx/make_index_question_delta.sh